# docker/supervisord.conf
# Continuous-feed process group for one demo tenant. The simulator-loop connects to
# the agent on 127.0.0.1:2407 (same container/network namespace), so all three must
# live together here. The loop tolerates the agent being momentarily down (it counts
# failed sends rather than crashing), so start order is not critical; autorestart
# heals anything that does die.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:agent]
command=php /app/artisan nightowl:agent
autostart=true
autorestart=true
startsecs=3
priority=10
stopwaitsecs=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:simulator-loop]
command=php /app/artisan nightowl:simulator-loop --token=%(ENV_NIGHTOWL_TOKEN)s
autostart=true
autorestart=true
startsecs=5
priority=20
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Drives nightowl:prune (scheduled when NIGHTOWL_SIMULATOR_ENABLED=true) so each
# tenant DB's raw telemetry stays bounded (~6h); charts keep filling from rollups.
[program:scheduler]
command=php /app/artisan schedule:work
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
