# docker/supervisord.conf
# Continuous-feed process group for one demo tenant. The simulator-loop connects to
# the agent on 127.0.0.1:2407 (same container/network namespace), so all three must
# live together here. The loop tolerates the agent being momentarily down (it counts
# failed sends rather than crashing), so start order is not critical; autorestart
# heals anything that does die.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:agent]
command=php /app/artisan nightowl:agent
autostart=true
autorestart=true
startsecs=3
priority=10
stopwaitsecs=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:simulator-loop]
command=php /app/artisan nightowl:simulator-loop --token=%(ENV_NIGHTOWL_TOKEN)s
autostart=true
autorestart=true
startsecs=5
priority=20
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Prune each tenant DB's raw telemetry to ~6h on a 30-min loop so the DBs stay bounded
# (charts keep filling from rollups). Runs nightowl:prune DIRECTLY rather than via
# `schedule:work`, whose ScheduleRunCommand output layer throws on a missing `task`
# console component under this stack — that crash-loop is what flips the container
# unhealthy. The direct loop never touches that path.
[program:prune]
command=/bin/sh -c "while true; do sleep 1800; php /app/artisan nightowl:prune --hours=6 || true; done"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
